# deploy/trainspotter.service
[Unit]
Description=Trainspotter Telegram Bot — Italo Treno price checker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USERNAME
WorkingDirectory=HOME_DIR/trainspotter
Environment=PYTHONPATH=HOME_DIR/trainspotter/src
ExecStart=HOME_DIR/trainspotter/.venv/bin/python -m trainspotter.bot
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
